# extras/upstart.conf
description "AMQProxy"
start on started networking
stop on shutdown
respawn
setuid nobody
setgid nogroup
console log
exec /usr/local/bin/amqproxy --listen=0.0.0.0 --port 5673 --default-prefetch=100
